# deploy/orchestrator-dashboard.service
[Unit]
Description=HyperOrchestrator Dashboard (Next.js)
After=network.target orchestrator-api.service
Wants=orchestrator-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/orch-cloud/dashboard
EnvironmentFile=-/opt/orch-cloud/.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=INTERNAL_API_URL=http://127.0.0.1:8000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
